# systemd/linux/wisp.service
[Unit]
Description=Wisp
After=network-online.target libvirtd.service virtqemud.service containerd.service
Wants=network-online.target

[Service]
User=WISP_USER
WorkingDirectory=WISP_PATH/backend
EnvironmentFile=-WISP_PATH/config/runtime.env
# Ensure the mDNS stub IP (169.254.53.53) is on br0 before the in-process DNS
# forwarder (backend/src/lib/linux/mdnsForwarder.js) tries to bind to it. The
# address is runtime-only — dropped on every reboot — so we re-add it here
# even though container-dns.sh also adds it at setup time. `+` runs this as
# root even though the main service runs as WISP_USER.
ExecStartPre=+/bin/sh -c '[ -d /sys/class/net/br0/bridge ] && (ip -4 addr show dev br0 | grep -q "169.254.53.53/32" || ip addr add 169.254.53.53/32 dev br0 scope link) || true'
ExecStart=/usr/bin/node src/index.js
# Grant CAP_NET_BIND_SERVICE so mdnsForwarder can bind UDP+TCP 169.254.53.53:53
# while the service otherwise runs unprivileged as WISP_USER. Do NOT also set
# CapabilityBoundingSet — it would cap every child process (including sudo)
# to just this one capability, and sudo needs CAP_SETUID/CAP_SETGID to elevate,
# which breaks the wisp-netns and wisp-mount privileged helpers.
AmbientCapabilities=CAP_NET_BIND_SERVICE
Restart=on-failure
RestartSec=5
TimeoutStopSec=15

[Install]
WantedBy=multi-user.target
